# config/chaos-consolidator.service.template
[Unit]
Description=CHAOS Memory Consolidator (Auto-Capture)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=%USER%
WorkingDirectory=%CHAOS_HOME%
Environment="HOME=%HOME%"
Environment="CHAOS_HOME=%CHAOS_HOME%"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:%CHAOS_HOME%/bin"

# Pre-flight checks (Ollama only - embedded Dolt auto-starts)
ExecStartPre=/bin/bash -c 'timeout 5 curl -sf http://localhost:11434/api/tags >/dev/null || (echo "Ollama not ready" && exit 1)'

# Main process
ExecStart=%CHAOS_HOME%/bin/chaos-consolidator --config %CHAOS_HOME%/config/consolidator.yaml --auto-capture

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chaos-consolidator

# Resource limits
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
